# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ── FastAPI backend (uvicorn) ─────────────────────────────────────────────────
[program:backend]
command=python -m uvicorn backend:app --host 0.0.0.0 --port 8000 --workers 1 --log-level info
directory=/app
autostart=true
autorestart=true
startretries=5
stderr_logfile=/var/log/supervisor/backend.err.log
stdout_logfile=/var/log/supervisor/backend.out.log
environment=PYTHONUNBUFFERED="1"

# ── Nginx (serves frontend + proxies /api to backend) ────────────────────────
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log

# ── Coordinator agent ─────────────────────────────────────────────────────────
[program:coordinator]
command=python agents/coordinator.py
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=5
stderr_logfile=/var/log/supervisor/coordinator.err.log
stdout_logfile=/var/log/supervisor/coordinator.out.log
environment=PYTHONUNBUFFERED="1"

# ── Event Intelligence agent ─────────────────────────────────────────────────
[program:event_intelligence]
command=python agents/event_intelligence.py
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=8
stderr_logfile=/var/log/supervisor/event_intelligence.err.log
stdout_logfile=/var/log/supervisor/event_intelligence.out.log
environment=PYTHONUNBUFFERED="1"

# ── Supplier Impact agent ─────────────────────────────────────────────────────
[program:supplier_impact]
command=python agents/supplier_impact.py
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=10
stderr_logfile=/var/log/supervisor/supplier_impact.err.log
stdout_logfile=/var/log/supervisor/supplier_impact.out.log
environment=PYTHONUNBUFFERED="1"

# ── Financial Exposure agent ──────────────────────────────────────────────────
[program:financial_exposure]
command=python agents/financial_exposure.py
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=12
stderr_logfile=/var/log/supervisor/financial_exposure.err.log
stdout_logfile=/var/log/supervisor/financial_exposure.out.log
environment=PYTHONUNBUFFERED="1"

# ── Regulatory & Trade agent ──────────────────────────────────────────────────
[program:regulatory_trade]
command=python agents/regulatory_trade.py
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=12
stderr_logfile=/var/log/supervisor/regulatory_trade.err.log
stdout_logfile=/var/log/supervisor/regulatory_trade.out.log
environment=PYTHONUNBUFFERED="1"

# ── Alt Sourcing agent ────────────────────────────────────────────────────────
[program:alt_sourcing]
command=python agents/alt_sourcing.py
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=14
stderr_logfile=/var/log/supervisor/alt_sourcing.err.log
stdout_logfile=/var/log/supervisor/alt_sourcing.out.log
environment=PYTHONUNBUFFERED="1"
